# expenses-web.service
[Unit]
Description=Expenses Web Dashboard
After=network.target

[Service]
Type=simple
User=blackhawk
WorkingDirectory=/home/blackhawk/dev/expenses-web
Environment=FLASK_PORT=6000
Environment=FLASK_DEBUG=false
ExecStart=/home/blackhawk/dev/expenses-web/venv/bin/python server.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
